# .config/systemd/user/test-timer.timer
[Unit]
Description=Test Timer

[Timer]
OnCalendar=*:*:0/20

[Install]
WantedBy=timers.target
