# .config/systemd/user/test-timer.service
[Unit]
Description=Test Timer Service

[Service]
Type=oneshot
ExecStart=bash -c "echo 'test' >> /tmp/test.service.log"
