# .config/systemd/user/xsession.target
[Unit]
Description=XSession Running
BindsTo=graphical-session.target
